# platform/scripts/launchagents-canonical/com.gathering.edge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- #3611 UNTANGLE — replaces com.chorus.caddy. :3000 is gathering's front
         door, so the service binding it lives in the gathering family: it stays
         up when the chorus stack stops (invariant #6, test-product-membrane.sh).
         Config at the runtime home ~/.gathering/data/caddy/Caddyfile (canonical
         source: gathering-edge-Caddyfile beside this file), never a repo path. -->
    <key>Label</key>
    <string>com.gathering.edge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/caddy/bin/caddy</string>
        <string>run</string>
        <string>--config</string>
        <string>/Users/jeffbridwell/.gathering/data/caddy/Caddyfile</string>
        <string>--adapter</string>
        <string>caddyfile</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jeffbridwell/.gathering/data/caddy</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/edge-caddy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/edge-caddy.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>XDG_DATA_HOME</key>
        <string>/opt/homebrew/var/lib</string>
        <key>HOME</key>
        <string>/Users/jeffbridwell</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
