# config/launchagents/com.chorus.reindex-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chorus.reindex-worker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/chorus-reindex-worker.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/.chorus/reindex-worker.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/.chorus/reindex-worker.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <!-- #3086: nvm v20.11.1 first so PATH-resolved tools match chorus-api's node.
             Belt-and-suspenders — chorus-reindex-worker.sh now pins node by absolute
             path, but keeping committed == the working installed config avoids the
             divergence that hid the #3085 ABI crash. -->
        <string>/Users/jeffbridwell/.nvm/versions/node/v20.11.1/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
