# roles/silas/scripts/launchagents/com.gathering.loki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gathering.loki</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jeffbridwell/bin/loki</string>
        <string>-config.file=/Users/jeffbridwell/CascadeProjects/shared-observability/config/loki/loki-native.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/loki.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/loki.log</string>
</dict>
</plist>
