# proving/config/launchagents/com.gathering.context-cache-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gathering.context-cache-daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
SCRIPT=/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/context-cache-daily.sh
for role in wren silas kade; do "$SCRIPT" "$role" 2>/dev/null; done
        </string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/context-cache-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/context-cache-daily-err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
