# proving/config/launchagents/com.chorus.session-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chorus.session-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeffbridwell/.chorus/scripts/chorus-watch-sessions.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/.chorus/watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/.chorus/watcher.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>CHORUS_ROOT</key>
        <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
    </dict>
</dict>
</plist>
