# proving/config/launchagents/com.chorus.hooks.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chorus.hooks</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/services/chorus-hooks/target/release/chorus-hooks</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/chorus-hooks.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/chorus-hooks.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>chorus_hooks=info</string>
        <key>CHORUS_ROOT</key>
        <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
    </dict>
</dict>
</plist>
