# proving/config/launchagents/com.chorus.daily-review-quality.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chorus.daily-review-quality</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/daily-review-quality.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>3</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/daily-review-quality.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/daily-review-quality.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <!-- .cargo/bin is required: daily-review-quality.sh re-runs the cargo suites,
             and without it `cargo` is not found (rc=127) → every crate is misreported
             as "1/1 failed" (compile/run failure). Matches nightly-suites.plist's PATH
             so the reporter's run and the 03:00 run resolve the same toolchain. (#3254) -->
        <string>/Users/jeffbridwell/.nvm/versions/node/v20.11.1/bin:/Users/jeffbridwell/.local/bin:/Users/jeffbridwell/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/jeffbridwell</string>
        <key>CHORUS_ROOT</key>
        <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
    </dict>
</dict>
</plist>
