# proving/config/launchagents/com.gathering.fuseki.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gathering.fuseki</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/fuseki-server</string>
        <string>--config=/Users/jeffbridwell/CascadeProjects/jeff-bridwell-personal-site/data/fuseki/configuration/pods-native.ttl</string>
        <string>--port=3030</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>JAVA_HOME</key>
        <string>/opt/homebrew/opt/openjdk</string>
        <key>JVM_ARGS</key>
        <string>-Xmx1536m</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/fuseki.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/fuseki.log</string>
</dict>
</plist>
