# proving/config/launchagents/com.gathering.log-rotate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gathering.log-rotate</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
LOG_DIR=/Users/jeffbridwell/Library/Logs/Gathering
for LOG in fuseki.log gathering-app.log loki.log prometheus.log grafana.log; do
  FILE="$LOG_DIR/$LOG"
  if [ -f "$FILE" ] &amp;&amp; [ $(stat -f%z "$FILE" 2>/dev/null || echo 0) -gt 104857600 ]; then
    mv "$FILE" "${FILE}.$(date +%Y%m%d)"
    gzip "${FILE}.$(date +%Y%m%d)" 2>/dev/null
    find "$LOG_DIR" -name "${LOG}.*" -mtime +7 -delete 2>/dev/null
  fi
done
        </string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/log-rotate.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jeffbridwell/Library/Logs/Gathering/log-rotate.log</string>
</dict>
</plist>
