# proving/config/launchagents/com.chorus.context-cache-5min.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chorus.context-cache-5min</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/context-cache-5min.sh silas; /Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/context-cache-5min.sh kade; /Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/context-cache-5min.sh wren</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/context-cache-5min.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/context-cache-5min.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/jeffbridwell</string>
        <key>CHORUS_ROOT</key>
        <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
    </dict>
</dict>
</plist>
