# proving/config/launchagents/com.chorus.fuseki-compact.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chorus.fuseki-compact</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/roles/silas/scripts/fuseki-maintenance.sh</string>
    <string>compact</string>
  </array>

  <!-- Saturday at 1:00 AM -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>6</integer>
    <key>Hour</key>
    <integer>1</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/fuseki-compact.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/fuseki-compact.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CHORUS_ROOT</key>
    <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
  </dict>
</dict>
</plist>
