# proving/config/launchagents/com.chorus.posture-capture.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chorus.posture-capture</string>

  <key>ProgramArguments</key>
  <array>
    <string>open</string>
    <string>-a</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/apps/PostureCapture.app</string>
    <string>--args</string>
    <string>--background</string>
  </array>

  <!-- Run every 5 minutes -->
  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/posture-capture.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/posture-capture.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CHORUS_ROOT</key>
    <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
  </dict>
</dict>
</plist>
