# proving/config/launchagents/com.chorus.launchagent-metrics.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.chorus.launchagent-metrics</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/jeffbridwell/CascadeProjects/shared-observability/scripts/launchagent-metrics.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>15</integer>
	<key>StandardErrorPath</key>
	<string>/tmp/launchagent-metrics.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/launchagent-metrics.log</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>CHORUS_ROOT</key>
		<string>/Users/jeffbridwell/CascadeProjects/chorus</string>
	</dict>
</dict>
</plist>
