# proving/config/launchagents/com.chorus.cruft-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chorus.cruft-scan</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/cruft-scan.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>259200</integer>
    <key>StandardOutPath</key>
    <string>/tmp/cruft-scan.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cruft-scan.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
