# platform/launchd/com.chorus.building-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.chorus.building-pipeline — fires building-pipeline.yml on release events.

  #2863: trigger is launchctl kickstart from chorus_acp's release-trigger
  step. /acp is the release event; this plist runs building-pipeline.yml
  whose first step (/build) fast-forwards canonical from origin/main, then
  builds and deploys. No WatchPaths, no scheduled poll — fires only on
  release.

  Pre-#2863 used WatchPaths on canonical's .git/refs/heads/main, paired
  with a chorus-werk-sync 10-min poll that kept canonical current. That
  introduced up to a 10-min staleness window after every peer commit and
  a separate failure surface (canonical-sync-aborted alert). Both retired
  by #2863. Sync now lives inside /build as an invariant first step.

  Install:
    cp platform/launchd/com.chorus.building-pipeline.plist ~/Library/LaunchAgents/
    launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.chorus.building-pipeline.plist

  Stop:
    launchctl bootout gui/$(id -u)/com.chorus.building-pipeline
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chorus.building-pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>cd /Users/jeffbridwell/CascadeProjects/chorus &amp;&amp; /opt/homebrew/bin/act push -P macos-latest=-self-hosted -W .github/workflows/building-pipeline.yml</string>
  </array>

  <key>StandardOutPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/building-pipeline.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/building-pipeline.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/jeffbridwell/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>CHORUS_ROOT</key>
    <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
    <key>DEPLOY_ROLE</key>
    <string>system</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
