# platform/launchd/com.chorus.crawler-index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.chorus.crawler-index — chorus crawler hydrator (#2817).

  Trigger: StartInterval=1800 ONLY — the crawler runs every 30 min.

  WatchPaths was REMOVED (#3068, Jeff's call 2026-05-24). It had fired the full
  27-domain sweep on every filesystem change under chorus/{designing,platform,
  roles,skills,directing,proving}, so commit/sync/deploy churn ran it constantly
  — a continuous eventloop-block storm. "Every 30 min is good enough; every 60
  seconds is a flood." StartInterval is now the SOLE trigger, not a fallback.
  ThrottleInterval was also removed — it debounced WatchPaths events; with no
  event triggers it serves no purpose. Same script (index-crawler-snapshots.sh);
  only the trigger surface changed.

  StandardOut/Err unify on a single .log file (already in Loki via
  Promtail scrape per #2824).

  Install:
    cp platform/launchd/com.chorus.crawler-index.plist ~/Library/LaunchAgents/
    launchctl unload -w ~/Library/LaunchAgents/com.chorus.crawler-index.plist  # if loaded with old plist
    launchctl load -w ~/Library/LaunchAgents/com.chorus.crawler-index.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chorus.crawler-index</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/index-crawler-snapshots.sh</string>
  </array>

  <!-- WatchPaths REMOVED (#3068, Jeff's call 2026-05-24). Firing the full
       27-domain sweep on every filesystem change under chorus/{designing,
       platform,roles,skills,directing,proving} turned commit/sync/deploy churn
       into a continuous eventloop-block STORM (1.7-12s blocks, ~27/10min).
       Empirical: with the crawler off, blocks dropped 6-11s -> ~1.3s. This is a
       FREQUENCY change, not a turn-off — every output is still produced on
       StartInterval below, so no consumer is starved. (Removing the crawler
       entirely would need the consumer dep-map #3071; dropping WatchPaths does
       not.) -->

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>StandardErrorPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/crawler-index.log</string>
  <key>StandardOutPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/crawler-index.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/jeffbridwell/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/jeffbridwell</string>
    <key>CHORUS_ROOT</key>
    <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
  </dict>
</dict>
</plist>
