# platform/launchd/com.chorus.crawler-index.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  com.chorus.crawler-index — chorus crawler hydrator (#2817).

  Replaces the prior StartCalendarInterval=15min calendar with file-watch
  primary + polling fallback. Same script as before (index-crawler-
  snapshots.sh); the trigger surface is what changes.

  Three trigger paths:
    1. WatchPaths fires on filesystem changes under chorus/{designing,
       platform, roles, skills, directing, proving} — sub-second after
       any role's edit.
    2. StartInterval=1800 polls every 30 min as a FALLBACK if WatchPaths
       events are missed (launchd event coalescing, host sleep, etc.).
       Throttled 60s->1800s (#3054, Silas): the 60s fallback drove ~200
       crawls/hr that queued behind each other on the spine. WatchPaths
       still fires sub-second on real edits, so this only widens the fallback.
    3. ThrottleInterval=5 enforces the debounce — no matter how many
       events fire in a 5-second window, launchd only invokes the
       script once per 5s. A `git pull` of N files triggers one run,
       not N.

  StandardOut/Err unify on a single .log file (already in Loki via
  Promtail scrape per #2824).

  Install:
    cp platform/launchd/com.chorus.crawler-index.plist ~/Library/LaunchAgents/
    launchctl unload -w ~/Library/LaunchAgents/com.chorus.crawler-index.plist  # if loaded with old plist
    launchctl load -w ~/Library/LaunchAgents/com.chorus.crawler-index.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chorus.crawler-index</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/index-crawler-snapshots.sh</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/designing</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/platform</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/roles</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/skills</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/directing</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/proving</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardErrorPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/crawler-index.log</string>
  <key>StandardOutPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/crawler-index.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/jeffbridwell/.cargo/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/jeffbridwell</string>
    <key>CHORUS_ROOT</key>
    <string>/Users/jeffbridwell/CascadeProjects/chorus</string>
  </dict>
</dict>
</plist>
