# platform/launchd/com.chorus.graph-hydrate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.chorus.graph-hydrate</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jeffbridwell/CascadeProjects/chorus/platform/scripts/graph-hydrate-tag.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DEPLOY_ROLE</key>
    <string>system</string>
  </dict>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/graph-hydrate.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jeffbridwell/Library/Logs/Chorus/graph-hydrate.log</string>
</dict>
</plist>
